# nfancurve.service
[Unit]
Description=Nfancurve service
After=default.target

[Service]
ExecStart=/bin/sh /usr/bin/nfancurve -c /etc/nfancurve.conf
KillSignal=SIGINT
Environment=LC_ALL=en_US.UTF-8

[Install]
WantedBy=default.target
